# crontab
* * * * * /var/www/trains/yii cron/delete

